# src/corewatcher.service.in
[Unit]
Description=Nitra corefile monitor for collecting crash data
DefaultDependencies=no
After=sysinit.target

[Service]
Type=forking
ExecStart=@prefix@/bin/corewatcher
Restart=always
StartLimitInterval=300
StartLimitBurst=2

[Install]
WantedBy=multi-user.target
